# code/marce-agent.service
# ============================================================
# AgentMarce - systemd service unit
# Install to: /etc/systemd/system/marce-agent.service
#
# Management commands:
#   sudo systemctl daemon-reload
#   sudo systemctl enable marce-agent   (auto-start on boot)
#   sudo systemctl start marce-agent
#   sudo systemctl status marce-agent
#   sudo journalctl -u marce-agent -f   (follow live logs)
# ============================================================

[Unit]
Description=AgentMarce - Homelab Orchestration Agent
Documentation=https://github.com/TrasteCC/AgentMarce
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=agentuser
Group=agentuser
WorkingDirectory=/home/agentuser/agent-services/antigravity-agent

# Load environment variables from .env
EnvironmentFile=/home/agentuser/agent-services/antigravity-agent/.env

# Use the virtual environment's Python binary
ExecStart=/home/agentuser/agent-services/antigravity-agent/.venv/bin/python api_server.py

# Restart automatically if the process crashes
Restart=always
RestartSec=10

# Resource limits (system protection)
MemoryMax=2G
CPUQuota=150%

# Service hardening
NoNewPrivileges=yes
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=marce-agent

[Install]
WantedBy=multi-user.target
